# configs/systemd/mnt-synology.automount
[Unit]
Description=Automount Synology NAS folder at boot

[Automount]
Where=/mnt/synology

[Install]
WantedBy=multi-user.target
