# configs/systemd/mnt-nas.automount
[Unit]
Description=Automount Synology NAS folder at boot

[Automount]
Where=/mnt/nas

[Install]
WantedBy=multi-user.target
